# src/templates/containerd.service.tmpl
[Unit]
Description=Containerd MicroVM Service

[Service]
ExecStart={{ .BinPath }} 

[Install]
WantedBy=multi-user.target