# src/templates/flintlockd.service
[Unit]
Description=Flintlock Daemon Service

[Service]
ExecStart=/usr/sbin/foo-daemon

[Install]
WantedBy=multi-user.target