# src/templates/flintlockd.service.tmpl
[Unit]
Description=Flintlock Daemon Service

[Service]
ExecStart={{ .BinPath }}

[Install]
WantedBy=multi-user.target